# O_ft_server/ft_server/srcs/supervisord.conf
[supervisord]
nodaemon=ture

[program:nginx]
command=/usr/sbin/service nginx start
